# nfs-automount/nfs-automount.timer
[Unit]
Description=NFS Automount Timer

[Timer]
OnCalendar=minutely
AccuracySec=1m
Persistent=true

[Install]
WantedBy=default.target
