# clean-pkg/cleanpkg.timer
[Unit]
Description=Clean pacman cache packages timer

[Timer]
OnCalendar=weekly

[Install]
WantedBy=default.target
