# clean-pkg/cleanpkg.service
[Unit]
Description=Clean pacman cache packages service

[Service]
ExecStart=/usr/bin/paccache -r
