# flexget/flexget.timer
[Unit]
Description=flexget

[Timer]
OnCalendar=hourly

[Install]
WantedBy=default.target
