# mpd/mpd.service
[Unit]
Description=Mpd Service

[Service]
ExecStart=/usr/bin/mpd %h/.mpd/mpd.conf

[Install]
WantedBy=default.target
