# mpd/mpd.service
[Unit]
Description=Mpd User Service by %u

[Service]
ExecStart=/usr/bin/mpd %h/.mpd/mpd.conf

[Install]
WantedBy=default.target
