# syncthing/syncthing.service
[Unit]
Description=Syncthing
After=network.target

[Service]
ExecStart=/usr/bin/syncthing

[Install]
WantedBy=default.target
